# config/vehicle.rules
ATTRS{idVendor}=="0483", ACTION=="add", RUN+="/sbin/ip link set up can0 type can bitrate 500000"
